# pterodactyl/containers/pterodactyl-redis.container
[Unit]
Description=Redis container for Pterodactyl

[Container]
ContainerName=pterodactyl-redis
Image=docker.io/library/redis:alpine
Pod=pterodactyl.pod
Secret=pterodactyl_redis,type=env,target=REDIS_PASSWORD
Exec=/bin/sh -c "redis-server --requirepass $${REDIS_PASSWORD}"

[Service]
Restart=unless-stopped

[Install]
WantedBy=default.target